# svconf/redis.conf
[program:redisd]
command=/usr/bin/redis-server
stdout_logfile=/home/vanshin/log/redisd.log
stdout_logfile_maxbytes=16MB
stdout_logfile_backups=10
